# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/crond -n
autostart=true
autorestart=true

[program:minecraft]
command=/bin/bash /opt/minecraft/run/systemd/execstart.sh
directory=/opt/minecraft/run
autostart=true
autorestart=true
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stopwaitsecs=120
environment=SERVER_NAME="hasturian"
